# sudo nano /etc/systemd/system/rcontrol-listener.service
[Unit]
Description=RControl Listener
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/listener_module.py
Restart=always

[Install]
WantedBy=multi-user.target